# mysql/supervisord-mysqld.conf
[program:mysqld]
command=/mysql_start.sh
numprocs=1
autostart=true
autorestart=true